# package/system/reglinux-scripts/rules/80-usbgadget.rules
# Trigger soft-connect when cable state changes (if supported by the UDC driver)
SUBSYSTEMS=="extcon", ACTION=="change", RUN+="/usr/bin/usbgadget trigger"
