# package/reglinux/reglinux-scripts/rules/80-switch-screen.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/bash /usr/bin/system-switch-screen-checker-delayed"
