# board/batocera/fsoverlay/etc/profile.d/20-alias.sh
# ---- ALIAS VALUES ----
alias mc='mc -xc'
alias reglinux-check-updates='batocera-es-swissknife --update'
alias which='command -v'
