# board/fsoverlay/etc/profile.d/20-prompt.sh
export HOSTNAME="$(hostname)"
export PS1='${USER}@${HOSTNAME} $(pwd)$ '
