# board/fsoverlay/etc/profile.d/20-prompt.sh
export HOSTNAME="$(hostname)"
if [ -n "${PS1-}" ]; then
    PS1='\u@\h \w$ '
fi
