# board/rockchip/rk3128/fsoverlay/etc/udev/rules.d/60-drm.rules
SUBSYSTEM=="drm", ACTION=="add|change", RUN+="/usr/bin/regmsg setOutput auto"
